# deploy/targets/run/supervisord.d/mongodb.ini
[program:mongodb]
command=mongod
    --dbpath /data/db
    --logpath /data/log/mongodb.log
    --logappend
    --bind_ip '127.0.0.1'
    --journal
    --smallfiles
